# templates/backup.service
[Unit]
Description=Backup Service

[Service]
Type=oneshot
ExecStart=/opt/backup/backup
User={{ s3_backup_user }}
Group={{ s3_backup_group }}
